# app-generation/squad-orchestrator/supervisord.conf
# Supervisord Configuration for SuperCore Squad Orchestrator
# This file manages all critical processes for the orchestration system

[unix_http_server]
file=/tmp/supervisor-squad-orchestrator.sock   ; path to your socket file

[supervisord]
logfile=/Users/jose.silva.lb/LBPay/supercore/scripts/squad-orchestrator/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord-squad-orchestrator.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor-squad-orchestrator.sock

# ============================================================================
# Redis (Message Broker + Results Backend)
# ============================================================================
[program:redis]
command=redis-server --port 6379
directory=/Users/jose.silva.lb/LBPay/supercore/scripts/squad-orchestrator
autostart=true
autorestart=true
startsecs=3
startretries=3
stdout_logfile=/Users/jose.silva.lb/LBPay/supercore/scripts/squad-orchestrator/logs/redis.log
stderr_logfile=/Users/jose.silva.lb/LBPay/supercore/scripts/squad-orchestrator/logs/redis-error.log
user=jose.silva.lb
priority=10

# ============================================================================
# Celery Worker - Cards Queue
# ============================================================================
[program:celery-worker-cards]
command=python3 -m celery -A celery_app worker --loglevel=info --concurrency=2 --queues=cards --hostname=cards-worker@%%h
directory=/Users/jose.silva.lb/LBPay/supercore/scripts/squad-orchestrator
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=30
killasgroup=true
stopasgroup=true
stdout_logfile=/Users/jose.silva.lb/LBPay/supercore/scripts/squad-orchestrator/logs/celery-cards.log
stderr_logfile=/Users/jose.silva.lb/LBPay/supercore/scripts/squad-orchestrator/logs/celery-cards-error.log
user=jose.silva.lb
priority=20

# ============================================================================
# Celery Worker - Maintenance Queue (cleanup tasks)
# ============================================================================
[program:celery-worker-maintenance]
command=python3 -m celery -A celery_app worker --loglevel=info --concurrency=1 --queues=maintenance --hostname=maintenance-worker@%%h
directory=/Users/jose.silva.lb/LBPay/supercore/scripts/squad-orchestrator
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=30
killasgroup=true
stopasgroup=true
stdout_logfile=/Users/jose.silva.lb/LBPay/supercore/scripts/squad-orchestrator/logs/celery-maintenance.log
stderr_logfile=/Users/jose.silva.lb/LBPay/supercore/scripts/squad-orchestrator/logs/celery-maintenance-error.log
user=jose.silva.lb
priority=20

# ============================================================================
# Portal Backend (FastAPI + WebSocket)
# ============================================================================
[program:portal-backend]
command=python3 server.py
directory=/Users/jose.silva.lb/LBPay/supercore/scripts/squad-orchestrator/monitoring/backend
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stdout_logfile=/Users/jose.silva.lb/LBPay/supercore/scripts/squad-orchestrator/logs/portal-backend.log
stderr_logfile=/Users/jose.silva.lb/LBPay/supercore/scripts/squad-orchestrator/logs/portal-backend-error.log
user=jose.silva.lb
priority=30

# ============================================================================
# Portal Frontend (npm dev server)
# ============================================================================
[program:portal-frontend]
command=npm run dev
directory=/Users/jose.silva.lb/LBPay/supercore/scripts/squad-orchestrator/monitoring/frontend
autostart=false
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=10
stdout_logfile=/Users/jose.silva.lb/LBPay/supercore/scripts/squad-orchestrator/logs/portal-frontend.log
stderr_logfile=/Users/jose.silva.lb/LBPay/supercore/scripts/squad-orchestrator/logs/portal-frontend-error.log
user=jose.silva.lb
priority=40

# ============================================================================
# Celery Beat (Periodic Tasks - Optional)
# ============================================================================
[program:celery-beat]
command=python3 -m celery -A celery_app beat --loglevel=info
directory=/Users/jose.silva.lb/LBPay/supercore/scripts/squad-orchestrator
autostart=false
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=10
stdout_logfile=/Users/jose.silva.lb/LBPay/supercore/scripts/squad-orchestrator/logs/celery-beat.log
stderr_logfile=/Users/jose.silva.lb/LBPay/supercore/scripts/squad-orchestrator/logs/celery-beat-error.log
user=jose.silva.lb
priority=25

# ============================================================================
# Group: All Essential Services
# ============================================================================
[group:squad-orchestrator]
programs=redis,celery-worker-cards,celery-worker-maintenance,portal-backend
priority=999
